# agents/forger/forger-v2.service
[Unit]
Description=Forger - Web Architect Agent v2.0 (with tools)
After=network.target agent-hub.service tool-bridge.service

[Service]
Type=simple
WorkingDirectory=/home/chad-yi/.openclaw/workspace/agents/forger
ExecStart=/usr/bin/python3 /home/chad-yi/.openclaw/workspace/agents/forger/forger-v2.py
Restart=always
RestartSec=60
Environment=PYTHONPATH=/home/chad-yi/.openclaw/workspace/infrastructure

[Install]
WantedBy=default.target
