# agents/_ARCHIVED_OLD_AGENTS/quanta-old-backup/quanta.service
[Unit]
Description=Quanta Trading Monitor v4.0
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/chad-yi/.openclaw/workspace/agents/quanta
ExecStart=/usr/bin/python3 /home/chad-yi/.openclaw/workspace/agents/quanta/monitor_callistofx.py
Restart=always
RestartSec=30

[Install]
WantedBy=default.target