# agents/quanta/quanta.service
[Unit]
Description=Quanta - Trading Bot (CALLISTOFX Monitor)
After=network.target

[Service]
Type=simple
User=chad-yi
WorkingDirectory=/home/chad-yi/.openclaw/workspace/agents/quanta
ExecStart=/usr/bin/python3 /home/chad-yi/.openclaw/workspace/agents/quanta/quanta_v4_service.py
Restart=always
RestartSec=60
Environment=PYTHONPATH=/home/chad-yi/.openclaw/workspace

[Install]
WantedBy=multi-user.target
