# agents/helios/helios.service
[Unit]
Description=Helios - Mission Control Auditor (COO)
After=network.target

[Service]
Type=simple
User=chad-yi
WorkingDirectory=/home/chad-yi/.openclaw/workspace/agents/helios
ExecStart=/usr/bin/python3 /home/chad-yi/.openclaw/workspace/agents/helios/helios-audit.py
Restart=always
RestartSec=900
Environment=PYTHONPATH=/home/chad-yi/.openclaw/workspace

[Install]
WantedBy=multi-user.target
